# web/systemd/nanobk-web-panel.service.example
[Unit]
Description=NanoBK Web Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace /opt/NanoBK-Proxy-Suite with your actual install path
WorkingDirectory=/opt/NanoBK-Proxy-Suite/web
ExecStart=/opt/NanoBK-Proxy-Suite/web/run.sh
Restart=on-failure
RestartSec=3
EnvironmentFile=/opt/NanoBK-Proxy-Suite/web/.env

# Security hardening
ProtectSystem=strict
# ReadWritePaths covers the whole web/ dir so first start can create/update .venv.
# The Web app still remains local-only and secrets stay in EnvironmentFile.
ReadWritePaths=/opt/NanoBK-Proxy-Suite/web
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
