# web/systemd/nanobk-web-panel.service.example
[Unit]
Description=NanoBK Web Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace /opt/NanoBK-Proxy-Suite with your actual install path
WorkingDirectory=/opt/NanoBK-Proxy-Suite/web
ExecStart=/opt/NanoBK-Proxy-Suite/web/run.sh
Restart=on-failure
RestartSec=3
EnvironmentFile=/opt/NanoBK-Proxy-Suite/web/.env

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/NanoBK-Proxy-Suite/web/.venv
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
